# tornado.conf
description "Used to start python script as a service"
author "bilbo"
start on vagrant-ready
stop on runlevel [06]
exec python /vagrant/chat.py
respawn